# ubuntucore/upstart/foreman-discovery.conf
start on (started firmware-update)
script
    /opt/foreman_startup.rb
    /opt/discovery_init.sh
end script
