# ubuntucore/upstart/getswport.conf
start on (started networking)
script
    cdpr -d eth0
end script
