# foreman_discovery/foreman-discovery.conf
start on (started networking)
script
    echo "exec /opt/foreman-discovery/foreman_startup.rb"
    echo "exec /opt/foreman-discovery/discovery_init.sh"
end script
