# foreman_discovery/foreman-discovery.conf
start on (started networking)
script
    /opt/foreman_startup.rb"
    /opt/discovery_init.sh"
end script
